# api-server/systemd/vk-api.service
[Unit]
Description=Visual Knowledge API Server
After=network.target
Wants=network-online.target
# Restart limits to prevent flapping
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/ubuntu/data-uploader

# Environment variables
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_ENV=production"

# Start command - logs to systemd journal AND file
ExecStart=/usr/bin/python3 main.py

# Restart configuration
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
# Memory limit - kill if exceeds 1GB
MemoryMax=1G
MemoryHigh=800M

# Send output to both systemd journal and log file
StandardOutput=append:/var/log/vk-api/app.log
StandardError=append:/var/log/vk-api/error.log

# Security settings (commented out for now as they might break functionality)
# PrivateTmp=yes
# NoNewPrivileges=true

[Install]
WantedBy=multi-user.target